# roles/software/templates/etc/profile.d/java_antialiasing.sh
#!/bin/sh
export _JAVA_OPTIONS="$_JAVA_OPTIONS -Dawt.useSystemAAFontSettings=lcd"
